# systemd/hostapd-wait.service
[Unit]
Description=Wait for uap0 interface before hostapd
Before=hostapd.service
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'TIMEOUT=30; ELAPSED=0; while [ $ELAPSED -lt $TIMEOUT ]; do if ip link show uap0 2>/dev/null | grep -q "state UP\|state UNKNOWN"; then exit 0; fi; sleep 1; ELAPSED=$((ELAPSED+1)); done; echo "Timeout waiting for uap0"; exit 1'
RemainAfterExit=yes

[Install]
WantedBy=hostapd.service

