# systemd/trichogramma.service
[Unit]
Description=Trichogramma Pi Bluetooth Service
Documentation=https://github.com/trichogramma-system
After=network.target bluetooth.target
Requires=bluetooth.service
# Garante que o serviço só inicia após o Bluetooth estar pronto
Wants=bluetooth.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/trichogramma-pi/service
ExecStart=/usr/bin/python3 /home/pi/trichogramma-pi/service/main.py

# Reinicia automaticamente em caso de falha
Restart=on-failure
RestartSec=5

# Limite de reinicializações (5 tentativas em 60 segundos)
StartLimitInterval=60
StartLimitBurst=5

# Logging via syslog/journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trichogramma

# Variáveis de ambiente (se necessário)
Environment="PYTHONUNBUFFERED=1"

# Segurança: previne que o serviço escreva em locais do sistema
# (Descomente se quiser restrições adicionais)
# ProtectSystem=strict
# ProtectHome=read-only
# ReadWritePaths=/var/log /home/pi/trichogramma-pi

# Timeout para inicialização
TimeoutStartSec=30

# Timeout para parada (shutdown gracioso)
TimeoutStopSec=10

# Sinal para enviar ao parar o serviço (SIGTERM para shutdown limpo)
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

