# systemd/uap0-setup.service
[Unit]
Description=Setup uap0 interface IP address
After=network.target
Before=hostapd.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'TIMEOUT=30; ELAPSED=0; while [ $ELAPSED -lt $TIMEOUT ]; do if ip link show uap0 2>/dev/null | grep -q "state UP\|state UNKNOWN"; then ip addr flush dev uap0 2>/dev/null; ip addr add 10.3.141.1/24 dev uap0; ip link set uap0 up; exit 0; fi; sleep 1; ELAPSED=$((ELAPSED+1)); done; echo "Timeout: uap0 interface not found"; exit 1'

[Install]
WantedBy=multi-user.target

